# forecast.timer
[Unit]
Description=Run Solar Forecast Once Per Day

[Timer]
OnCalendar=*-*-* 08:30:00
Persistent=true

[Install]
WantedBy=timers.target